# cpu-boost
#!/bin/sh /etc/rc.common

START=99

start() {

logger -t tweak-cpu "Starting Tweak Cpu Core"

echo 1 > /sys/kernel/debug/ecm/ecm_classifier_default/accel_delay_pkts
echo 1 > /sys/kernel/debug/ecm/ecm_db/defunct_all
echo f > /proc/net/nf_conntrack
echo 1 > /proc/sys/net/netfilter/nf_conntrack_events

}

stop() {

logger -t tweak cpu "Stopping Tweak Cpu Core"

echo 1 > /sys/kernel/debug/ecm/ecm_classifier_default/accel_delay_pkts
echo 1 > /sys/kernel/debug/ecm/ecm_db/defunct_all
echo f > /proc/net/nf_conntrack
echo 1 > /proc/sys/net/netfilter/nf_conntrack_events

}
