# cpu-boost
#!/bin/sh /etc/rc.common

START=99

start() {

logger -t tweak-cpu "Starting Tweak Cpu Core"

IFACE=$(ls /sys/class/net)
	
for i in ${IFACE}; do
	ethtool -K $i gro on 2>/dev/null
	if [[ -e /sys/class/net/$i/queues/rx-0/rps_cpus ]]; then
		if [[ $i = "wwan0" ]]; then
			echo f > /sys/class/net/$i/queues/rx-0/rps_cpus

		else
			echo f > /sys/class/net/$i/queues/rx-0/rps_cpus

		fi
	fi
done

}

stop() {

logger -t tweak cpu "Stopping Tweak Cpu Core"

IFACE=$(ls /sys/class/net)
	
for i in ${IFACE}; do
	ethtool -K $i gro on 2>/dev/null
	if [[ -e /sys/class/net/$i/queues/rx-0/rps_cpus ]]; then
		if [[ $i = "wwan0" ]]; then
			echo f > /sys/class/net/$i/queues/rx-0/rps_cpus

		else
			echo f > /sys/class/net/$i/queues/rx-0/rps_cpus

		fi
	fi
done

}
