# rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# Restart passwall after power on
/etc/init.d/passwall restart
# Running script custom ttl64
/etc/init.d/firewall-custom restart
/etc/init.d/cpu-boost restart

sysctl net.core.default_qdisc=fq_codel
syactl net.ipv4.tcp_congestion_control=bbr

exit 0

